# cloudflare-ddns/systemd-service/cloudflare-ddns.timer
[Unit]
Description=Update DDNS on Cloudflare every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target