# on-boot-script/examples/unifi-os-files/udmboot.service
[Unit]
Description=Run On Startup UDM
After=network.target

[Service]
ExecStart=/etc/init.d/udm.sh

[Install]
WantedBy=multi-user.target
